# deploy/systemd/br_tcpdump_rotate.service
[Unit]
Description=BlackRoad ledger rotation (one-shot)
After=network.target

[Service]
Type=oneshot
Environment=REPO_DIR=/opt/blackroad
ExecStart=/bin/bash -c "$REPO_DIR/scripts/br_rotate_ledger.sh"

[Install]
WantedBy=multi-user.target
