# deploy/systemd/br_tcpdump_rotate.timer
[Unit]
Description=Daily rotation timer for BlackRoad ledger

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
